# docker/supervisor.d/pwa.dev.ini
[program:pwa]
directory=/srv/pwa
command=pnpm dev --host
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/pwa.log
stdout_logfile_maxbytes=1MB
redirect_stderr=true
